# root/etc/s6-overlay/s6-rc.d/telegram-bot/run
#!/command/execlineb -P
with-contenv
s6-setuidgid abc
/app/bot
